# Gatesentry/rbi/rbi-chrome/chromium-kiosk.template.conf
; Rendered by entrypoint.sh (envsubst ${ALLOWED_URL}) -> /etc/neko/supervisord/chromium.conf
; A COMPLETE Chromium runs in the container, but the client view is KIOSK-LOCKED to
; the one allowed URL: no tabs, no address bar, no new-tab — can't navigate away.
[program:chromium]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s",VK_ICD_FILENAMES="/usr/share/vulkan/icd.d/lvp_icd.json",VK_DRIVER_FILES="/usr/share/vulkan/icd.d/lvp_icd.json"
command=/usr/bin/chromium --no-sandbox --window-position=0,0 --window-size=1920,1080 --force-device-scale-factor=1.25 --kiosk --display=%(ENV_DISPLAY)s --user-data-dir=/home/neko/.config/chromium --no-first-run --use-gl=angle --use-angle=vulkan --enable-features=Vulkan --ignore-gpu-blocklist --disable-dev-shm-usage --disable-file-system --disable-pinch --overscroll-history-navigation=0 --load-extension=/opt/rbi/ext --disable-extensions-except=/opt/rbi/ext --disable-component-extensions-with-background-pages --no-default-browser-check --autoplay-policy=no-user-gesture-required --use-fake-ui-for-media-stream --disable-session-crashed-bubble --hide-crash-restore-bubble --user-agent="Mozilla/5.0 (Windows NT 10.0; Win64; x64) AppleWebKit/537.36 (KHTML, like Gecko) Chrome/140.0.0.0 Safari/537.36 Edg/140.0.0.0" ${ALLOWED_URL}
stopsignal=INT
autorestart=true
priority=800
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/chromium.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true

[program:openbox]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s",DISPLAY="%(ENV_DISPLAY)s"
command=/usr/bin/openbox --config-file /etc/neko/openbox.xml
autorestart=true
priority=300
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/openbox.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
redirect_stderr=true
